# files/etc/systemd/system/rs-collector.service.d/20-stdout-err.conf
# Delete this file to log to systemd journal.
[Service]
StandardOutput=file:/var/log/rs-collector.log
StandardError=inherit